# .config/systemd/system/restart-networkmanager.timer
[Unit]
Description=Restart NetworkManager every 5 hours

[Timer]
OnBootSec=5min
OnUnitActiveSec=5h
Persistent=true

[Install]
WantedBy=timers.target
